# roles/bifrozt/templates/var/spool/cron/crontabs/root
# m h  dom mon dow   command
*/3  *   *    *    *    /usr/sbin/logrotate /etc/logrotate.d/iptables >/dev/null 2>&1
